# roles/interactive/templates/jupyterhub.service.j2
[Unit]
Description=Jupyterhub

[Service]
User=root
ExecStart={{ apps_dir }}/jupyter-{{ jupyter_version }}/jupyter.bash

[Install]
WantedBy=multi-user.target
